# crontab
SHELL=/bin/bash
0 0 * * * root cd /opt/code && /usr/local/bin/python3 -u get_all_fonts.py >> /var/log/cron.log 2>&1

